# shell/crontab/global.config.sh
#!/bin/sh
# 数据库配置信息
DB_HOST=localhost               # 数据库主机地址
DB_USERNAME=root		# 数据库用户名
DB_PASSWORD=""            # 数据库密码
DB_DATABASE=WeCenter            # 数据库名
DB_DATABASE_SY=sy_isinho_com    # 沈阳这边的数据库

# ssh设置信息
SSH_HOST=www.ekotlin.com           # 远程主机地址
SSH_USERNAME=root            # 远程主机用户名
#SSH_PASSWORD=password         # 远程主机密码
SSH_PORT=8433                  # 远程主机端口
SSH_STORE_DIR=/tmp           # 远程保存数据备份文件地址

LOCAL_STORE_DIR=/tmp         # 本地保存地址
